# tests/supervisord.conf
[unix_http_server]
file = /tmp/test_services_supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/test_services_supervisor.sock

[supervisord]
pidfile=/tmp/supervisord.pid
logfile=/tmp/supervisord.log

[program:dynalite]
command=dynalite --port 8000 --createTableMs 1 --deleteTableMs 1 --updateTableMs 1
autostart=true
stdout_logfile=/tmp/dynalite_stdout.log
stderr_logfile=/tmp/dynalite_stderr.log

[program:gcd]
command=/tmp/gcdrun
autostart=true
stopasgroup=true
stdout_logfile=/tmp/gcd_stdout.log
stderr_logfile=/tmp/gcd_stderr.log

[program:s3server]
command=s3server.sh
autostart=true
stopasgroup=true
stdout_logfile=/tmp/s3server_stdout.log
stderr_logfile=/tmp/s3server_stderr.log
